# systemd/example.timer
[Unit]

[Timer]
OnCalendar=*-*-* 6,17:00

[Install]
WantedBy=timers.target
